# bitable-annotator.service
[Unit]
Description=Feishu Bitable Annotator Service
After=network.target

[Service]
Type=simple
User=yuntong.li
WorkingDirectory=/zyt_users/yuntong.li/server/bitable-annotator
# 启动前杀掉旧进程，防止端口冲突（- 前缀表示忽略失败）
ExecStartPre=-/usr/bin/pkill -f "node.*bitable-annotator-server.mjs"
ExecStartPre=-/bin/sleep 1
ExecStart=/usr/bin/node bitable-annotator-server.mjs
ExecStop=/bin/kill -SIGTERM $MAINPID
# 崩溃后 3 秒自动重启
Restart=always
RestartSec=3
# 日志输出到 systemd journal
StandardOutput=journal
StandardError=journal
Environment=PORT=6868
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
